# scripts/tpm-keyring/files/tpm-keyring-unlock.service
[Unit]
Description=Unlock GNOME keyring from TPM-bound credential
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/tpm-keyring-unlock.sh

[Install]
WantedBy=default.target